# ops/vps/rumia-valhalla-api.service
[Unit]
Description=Rumia loopback route contract adapter (preflight)
Documentation=file:///opt/apps/rumia/current/docs/ops/rumia-map-unblock-preflight.md
Wants=network-online.target
After=network-online.target
ConditionPathExists=/opt/rumia-map/valhalla-venv/bin/python
ConditionPathExists=/opt/rumia-map/valhalla/valhalla.json
ConditionPathExists=/opt/rumia-map/valhalla/tiles

[Service]
Type=simple
User=rumia
Group=rumia
ExecStart=/opt/rumia-map/valhalla-venv/bin/python /opt/rumia-map/valhalla/rumia-valhalla-api.py --config /opt/rumia-map/valhalla/valhalla.json --host 127.0.0.1 --port 3012
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadOnlyPaths=/opt/rumia-map/valhalla
MemoryMax=2G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
